# rustdesk-server/rootfs/etc/s6-overlay/s6-rc.d/hbbs/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: RustDesk Server
# Runs the hbbs (ID / rendezvous) service
# ==============================================================================

declare relay
declare -a params

cd /config || bashio::exit.nok "Could not change working directory to /config"

params=()

# Encrypted-only mode forces clients to use the server's key.
if bashio::config.true 'encrypted_only'; then
    bashio::log.info "Encrypted-only mode is enabled (-k _)"
    params+=(-k "_")
fi

# Relay server(s) for hbbs to advertise to clients.
if bashio::config.has_value 'relay'; then
    relay=$(bashio::config 'relay')
    bashio::log.info "Relay for hbbs set to: ${relay}"
    params+=(-r "${relay}")
fi

bashio::log.info "Starting hbbs: $(/opt/rustdesk-server/hbbs --version)"
exec /opt/rustdesk-server/hbbs "${params[@]}"
